# files/etc/supervisor/conf.d/saslauthd.conf
[program:saslauthd]
#command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
command=/bin/bash -c "/usr/sbin/saslauthd -d -a ldap -c -m /var/spool/postfix/var/run/saslauthd -n 5"
autostart=true
autorestart=true
startretries=3

# Redirect output so we can see it using "docker logs"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

